# .deploy/manlab-init.service
[Unit]
Description=ManLab AutoStart Service

[Service]
Type=simple
User=root
ExecStart=/root/ManLab/bin/ManLab

[Install]
WantedBy=multi-user.target